# live-build/config/includes.chroot/etc/profile.d/debz-live.sh
# shellcheck shell=bash
export DEBZ_LIVE=1
