# live-build/config/includes.chroot/usr/lib/systemd/system/debz-live-darksite.service
[Unit]
Description=Debz Live Darksite Ensure
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/lib/debz-installer/firstboot/ensure-darksite.sh

[Install]
WantedBy=multi-user.target
