# live-build/config/includes.chroot/usr/lib/systemd/system/debz-live-ssh-init.service
[Unit]
Description=debz live SSH keypair initialisation
Documentation=man:debz-live-ssh-init(8)
After=network.target
ConditionPathExists=!/home/live/.ssh/id_ed25519
# Only run on live session, not on installed target
ConditionPathExists=/etc/debz/build-marker

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/debz-live-ssh-init
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
