# live-build/config/includes.chroot/usr/lib/systemd/system/debz-srv-snapshot.service
[Unit]
Description=Debz /srv Workload ZFS Snapshot
After=zfs-mount.service
Wants=zfs-mount.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/snapshot-create.sh srv
StandardOutput=journal
StandardError=journal
