# live-build/config/includes.chroot/usr/lib/systemd/system/debz-autobootstrap.service
[Unit]
Description=debz cluster autobootstrap and WireGuard reflector
After=salt-master.service wg-quick@wg1.service
Requires=salt-master.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/debz-autobootstrap
StandardOutput=journal
StandardError=journal
