# live-build/config/includes.chroot/usr/lib/systemd/system/debz-snapshot.service
[Unit]
Description=Debz ZFS Snapshot Service

[Service]
Type=oneshot
ExecStart=/usr/sbin/debz-snapshot
