# live-build/config/includes.chroot/usr/lib/systemd/system/debz-autobootstrap.timer
[Unit]
Description=debz autobootstrap timer — peer reflection + cluster init
After=salt-master.service

[Timer]
OnBootSec=60s
OnUnitActiveSec=30s
AccuracySec=5s

[Install]
WantedBy=timers.target
