# live-build/config/includes.chroot/usr/lib/systemd/system/debz-webui.service
[Unit]
Description=debz Web UI (installer + management frontend)
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/debz-webui --port 8080 --no-browser
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
WorkingDirectory=/usr/local/share/debz-webui

[Install]
WantedBy=multi-user.target
