# live-build/config/includes.chroot/usr/lib/systemd/system/debz-firstboot.service
[Unit]
Description=Debz First Boot Initialization
After=network-online.target zfs-mount.service
Wants=network-online.target
ConditionPathExists=!/var/lib/debz/firstboot-done

[Service]
Type=oneshot
ExecStart=/usr/sbin/debz-firstboot
ExecStartPost=/bin/mkdir -p /var/lib/debz
ExecStartPost=/bin/touch /var/lib/debz/firstboot-done
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
