# live-build/config/includes.chroot/etc/systemd/system/debz-snapshot.timer
[Unit]
Description=Debz ZFS Snapshot Timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
