# live-build/config/includes.chroot/etc/systemd/system/debz-srv-snapshot.service
[Unit]
Description=Debz /srv Workload Snapshot
After=zfs-mount.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/snapshot-create.sh srv
