# live-build/config/includes.chroot/etc/systemd/system/debz-autoinstall.service
[Unit]
Description=Debz Unattended Installer (DEBZ-SEED disk)
Documentation=https://github.com/your-org/debz
After=local-fs.target systemd-udev-settle.service
Wants=systemd-udev-settle.service
# Only run in the live environment, not on the installed system
ConditionPathExists=!/var/lib/debz-installer/install-done

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/debz-autoinstall
StandardOutput=journal+console
StandardError=journal+console
RemainAfterExit=yes
# Give udev time to settle and label disks before scanning
ExecStartPre=/bin/sleep 3

[Install]
WantedBy=multi-user.target
