# live-build/config/includes.chroot/etc/systemd/system/debz-package-holds.service
[Unit]
Description=Apply Debz platform package holds
ConditionPathExists=/usr/sbin/debz-apply-platform-holds

[Service]
Type=oneshot
ExecStart=/usr/sbin/debz-apply-platform-holds

[Install]
WantedBy=multi-user.target
