# live-build/config/includes.chroot/etc/systemd/system/debz-apt-snapshot.service
[Unit]
Description=Debz APT Pre-Transaction ZFS Snapshot
DefaultDependencies=no
Before=apt-daily.service dpkg.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/snapshot-create.sh apt-pre
RemainAfterExit=yes
