# live-build/config/includes.chroot/etc/systemd/system/debz-srv-snapshot.timer
[Unit]
Description=Debz /srv Snapshot Timer (every 15 min, keep 4)

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
AccuracySec=1min

[Install]
WantedBy=timers.target
